# wdb.server.service
[Unit]
Description=Wdb server
After=syslog.target

[Service]
ExecStart=/usr/bin/wdb.server.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
