# server/wdb.server.socket
[Socket]
ListenStream=1984

[Install]
WantedBy=sockets.target
