# server/wdb.server.service
[Unit]
Description=Wdb server

[Service]
ExecStart=/usr/bin/wdb.server.py
