# server/wdb.server.socket
[Socket]
ListenStream=1984
ListenStream=19840

[Install]
WantedBy=sockets.target
